# smf/clickhouse/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='clickhouse'>

  <service name='system/illumos/clickhouse' type='service' version='1'>
  <create_default_instance enabled='false' />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <method_context>
    <method_environment>
      <envvar name="LD_LIBRARY_PATH" value="/opt/ooce/pgsql-13/lib/amd64" />
    </method_environment>
  </method_context>
  <exec_method type='method' name='start'
    exec='ctrun -l child -o noorphan,regent /opt/oxide/clickhouse/clickhouse server --log-file /var/tmp/clickhouse-server.log --errorlog-file /var/tmp/clickhouse-server.errlog -- --path %{config/store} --listen_host %{config/listen_host} --http_port 8123 &amp;'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='listen_host' type='astring' value='unknown' />
    <propval name='store' type='astring' value='unknown' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Clickhouse</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Oxide's Metrics Database</loctext>
    </description>
  </template>
</service>

</service_bundle>
