# smf/zone-network-setup/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='zone-network-setup'>

<service name='oxide/zone-network-setup' type='service' version='1'>
  <create_default_instance enabled='true' />

  <!-- Run after the zone's networking stack is up. -->
  <dependency name='network' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/network:default' />
  </dependency>

  <!-- Run after the NDP daemon is online.. -->
  <dependency name='ndp' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/network/routing/ndp:default' />
  </dependency>

  <!-- The zone-setup binary is not ready to run until its initial properties
          have been set by the sled-agent, which happens after the
          `manifest-import` service is running.
  -->
  <dependency name='manifest-import' type='service' grouping='require_all' restart_on='none'>
    <service_fmri value='svc:/system/manifest-import:default' />
  </dependency>

  <exec_method type='method' name='start'
    exec='/opt/oxide/zone-setup-cli/bin/zone-setup common-networking -d %{config/datalink} -s %{config/static_addr} -g %{config/gateway}'
    timeout_seconds='0' />

  <exec_method type='method' name='stop' exec=':true' timeout_seconds='0' />

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='transient' />
  </property_group>

  <property_group name='config' type='application'>
    <propval name='datalink' type='astring' value='unknown' />
    <propval name='gateway' type='astring' value='' />
    <propval name='static_addr' type='astring' value='unknown' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Zone Network Setup</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Configures networking for control plane zones</loctext>
    </description>
  </template>
</service>

</service_bundle>
