# smf/chrony-setup/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='chrony-setup'>

<service name='oxide/chrony-setup' type='service' version='1'>
  <create_default_instance enabled='true' />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <exec_method type='method' name='start'
    exec='/opt/oxide/zone-setup-cli/bin/zone-setup chrony-setup -b %{config/boundary} -p ${config/boundary_pool} -s %{config/server} -a %{config/allow}'
    timeout_seconds='0'>
  <method_context security_flags="aslr">
  <method_credential user="root" group="root"
    privileges="basic,file_chown" />
  </method_context>
  </exec_method>

  
  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='transient' />
  </property_group>

  <property_group name="config" type="application">
    <propval name="boundary" type="boolean" value="false" />
    <propval name="boundary_pool" type="astring" value="" />
    <propval name="server" type="astring" value="" />
    <propval name="allow" type="astring" value="" />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Chrony Setup</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Configures chronyd for the NTP zone</loctext>
    </description>
  </template>
</service>

</service_bundle>
