# smf/nexus/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='nexus'>

<service name='oxide/nexus' type='service' version='1'>
  <create_default_instance enabled='true' />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <dependency name='network' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/network:default' />
  </dependency>

  <dependency name='routing' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/network/routing/ndp:default' />
  </dependency>

  <dependency name='zone_network_setup' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/oxide/zone-network-setup:default' />
  </dependency>

  <dependency name='opte_interface_setup' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/oxide/opte-interface-setup:default' />
  </dependency>

  <method_context>
    <method_environment>
      <envvar name="LD_LIBRARY_PATH" value="/opt/ooce/pgsql-13/lib/amd64" />
    </method_environment>
  </method_context>
  <exec_method type='method' name='start'
    exec='ctrun -l child -o noorphan,regent /opt/oxide/omicron-nexus/bin/nexus /var/svc/manifest/site/nexus/config.toml &amp;'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Controller</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Oxide Control Plane Instance</loctext>
    </description>
  </template>
</service>

</service_bundle>
