# smf/cockroachdb/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='cockroachdb'>

  <service name='oxide/cockroachdb' type='service' version='1'>
  <create_default_instance enabled='true' />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <dependency name='zone_network_setup' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/oxide/zone-network-setup:default' />
  </dependency>

  <exec_method type='method' name='start'
    exec='/opt/oxide/lib/svc/manifest/cockroachdb.sh'
    timeout_seconds='0'>
    <method_context working_directory="/">
      <method_environment>
        <envvar name="GOTRACEBACK" value="crash" />
      </method_environment>
    </method_context>
  </exec_method>
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='listen_addr' type='astring' value='unknown' />
    <propval name='http_addr' type='astring' value='unknown' />
    <propval name='store' type='astring' value='unknown' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>CockroachDB</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Oxide's Distributed Database</loctext>
    </description>
  </template>
</service>

</service_bundle>
