# smf/cockroachdb/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='cockroachdb'>

  <service name='system/illumos/cockroachdb' type='service' version='1'>
  <create_default_instance enabled='true' />
  <single_instance />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <method_context>
    <method_environment>
      <envvar name="LD_LIBRARY_PATH" value="/opt/ooce/pgsql-13/lib/amd64" />
    </method_environment>
  </method_context>
  <!--
    TODO(https://github.com/oxidecomputer/omicron/issues/727)
    Call "start", joining with other crdb instances
  -->
  <exec_method type='method' name='start'
    exec='/var/svc/manifest/site/cockroachdb/method_script.sh'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='datalink' type='astring' value='unknown' />
    <propval name='gateway' type='astring' value='unknown' />
    <propval name='listen_addr' type='astring' value='unknown' />
    <propval name='listen_port' type='astring' value='unknown' />
    <propval name='store' type='astring' value='unknown' />
    <propval name='join_addrs' type='astring' value='unknown' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>CockroachDB</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Oxide's Distributed Database</loctext>
    </description>
  </template>
</service>

</service_bundle>
