# smf/propolis-server/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='propolis-server'>

<service name='system/illumos/propolis-server' type='service' version='1'>
  <dependency name='network' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/network:default' />
  </dependency>

  <method_context>
    <method_environment>
      <envvar name="LD_LIBRARY_PATH" value="/opt/ooce/pgsql-13/lib/amd64" />
    </method_environment>
  </method_context>
  <exec_method type='method' name='start'
    exec='ctrun -l child -o noorphan,regent /opt/oxide/propolis-server/bin/propolis-server run /var/svc/manifest/site/propolis-server/config.toml %{config/server_addr} &amp;'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='server_addr' type='astring' value='unknown' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Propolis Server</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Hypervisor</loctext>
    </description>
  </template>
</service>

</service_bundle>
