# smf/sled-agent/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='sled-agent'>

<service name='oxide/sled-agent' type='service' version='1'>
  <create_default_instance enabled='true' />

  <!-- File systems should all be mounted before we begin: -->
  <dependency name='fs-root' grouping='require_all' restart_on='none'
    type='service'>
    <service_fmri value='svc:/system/filesystem/root' />
  </dependency>
  <dependency name='fs-minimal' grouping='require_all' restart_on='none'
    type='service'>
    <service_fmri value='svc:/system/filesystem/minimal' />
  </dependency>

  <!-- Run after initial device configuration has completed. -->
  <dependency name='devices' grouping='require_all' restart_on='none'
    type='service'>
    <service_fmri value='svc:/milestone/devices' />
  </dependency>

  <!-- The bootstrapping portion of the agent uses the network to communicate
    with other servers on the rack. -->
  <dependency name='network' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/network' />
  </dependency>

  <!-- ... and make sure we run after the T6 has been configured. This is
  grouped as `optional_all` to allow sled-agent to run on devices that do
  not have a t6init service (e.g., testing in CI). -->
  <dependency name='t6init' grouping='optional_all' restart_on='none'
    type='service'>
        <service_fmri value='svc:/system/t6init' />
  </dependency>

  <exec_method type='method' name='start'
    exec='ctrun -l child -o noorphan,regent /opt/oxide/sled-agent/sled-agent run /opt/oxide/sled-agent/pkg/config.toml &amp;'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Sled Agent</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Per-Server Sled Agent</loctext>
    </description>
  </template>
</service>

</service_bundle>
