# smf/ntp-admin/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='ntp-admin'>

  <service name='oxide/ntp-admin' type='service' version='1'>
  <create_default_instance enabled='true' />

  <dependency name='multi_user' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/multi-user:default' />
  </dependency>

  <dependency name='zone_network_setup' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/oxide/zone-network-setup:default' />
  </dependency>

  <dependency name='ntp' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/oxide/ntp:default' />
  </dependency>

  <exec_method type='method' name='start'
    exec='/opt/oxide/lib/svc/manifest/ntp-admin.sh'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='address' type='astring' value='unknown' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='contract' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Omicron NTP Admin</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Administering Oxide's NTP Service</loctext>
    </description>
  </template>
</service>

</service_bundle>
